# atomicmirror/atomicmirror_template.service
[Unit]
Description=atomicmirror NAME
After=Docker.service

[Service]
ExecStart=/usr/bin/docker start NAME
ExecStop=/usr/bin/docker stop NAME

[Install]
WantedBy=multi-user.target
